# systemd/live-boxscores.service
[Unit]
Description=Backfill NBA boxscores for the previous slate
After=network-online.target live-pipeline.service
Wants=network-online.target

[Service]
Type=oneshot
User=daniel
Group=daniel
WorkingDirectory=/home/daniel/projects/projections-v2
Environment=PROJECTIONS_DATA_ROOT=/home/daniel/projections-data
ExecStart=/home/daniel/projects/projections-v2/scripts/run_boxscores.sh
StandardOutput=journal
StandardError=journal
Nice=5

[Install]
WantedBy=multi-user.target
