# systemd/nightly-eval.service
[Unit]
Description=Run nightly prediction accuracy evaluation (FPTS + ownership)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/daniel/projects/projections-v2
Environment=PROJECTIONS_DATA_ROOT=/home/daniel/projections-data
User=daniel
Group=daniel
ExecStart=/home/daniel/.local/bin/uv run python -m scripts.analyze_accuracy --output /home/daniel/projections-data/reports/eval_latest.json
StandardOutput=journal
StandardError=journal
Nice=5

[Install]
WantedBy=multi-user.target
